# ops/kamiyo-agent-exec/kamiyo-agent-assessment-growth.timer
[Unit]
Description=Run Kamiyo Agent assessment/dispute growth daily

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=8h
AccuracySec=5min
Unit=kamiyo-agent-assessment-growth.service
Persistent=true

[Install]
WantedBy=timers.target
